# production/conf/supervisor.conf
[program:dj_docs]
command=/home/djdocs/dj_docs/production/gunicorn_start
user=root
stdout_logfile=/home/djdocs/dj_docs/log/supervisor.log
autorestart=true
redirect_stderr=true
