# conf/supervisor.conf
[program:dj_docs]
command=/home/djdocs/bin/gunicorn dj_docs.wsgi -c /home/djdocs/dj_docs/gunicorn.conf.py
directory=/home/djdocs/dj_docs
user=root
autorestart=true
redirect_stderr=true
